# docker-deps/supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=/var/www/sellerapp &
autorestart=true
